# enhance/rel.2.3/ibm/dist/etc/inittab
is:s:initdefault:
r0:0:wait:/etc/uadmin 2 0
